# web/deploy/oca.service
[Unit]
Description=Option Chain Archive (Next.js standalone)
After=network.target

[Service]
Type=simple
User=oca
Group=oca
WorkingDirectory=/opt/oca/web/.next/standalone
EnvironmentFile=/opt/oca/web/.env.production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
